# systemd/.config/systemd/user/layer-shell.service
[Unit]
Description=Layer Shell
After=pipewire-dbus.service

[Service]
ExecStart=layer-shell
WorkingDirectory=%h
Restart=on-failure
RestartSec=5s
