# systemd/.config/systemd/user/pipewire-dbus.service
[Unit]
Description=Pipewire DBus

[Service]
ExecStart=pipewire-dbus
WorkingDirectory=%h
Restart=on-failure
RestartSec=5s
