# systemd/.config/systemd/user/swaybg.service
[Unit]
Description=SwayBG
After=graphical-session.target

[Service]
ExecStart=swaybg --mode stretch --image .local/share/backgrounds/current.jpeg
WorkingDirectory=%h
Restart=on-failure
RestartSec=5s
Slice=app-graphical.slice

[Install]
WantedBy=graphical-session.target
